# etc/udev/rules.d/redGearCosmos71_autoeffects/86-disable-requireda.rules
#SUBSYSTEM=="net", SUBSYSTEMS=="pci", ATTRS{idVendor}=="10ec", RUN+="/usr/bin/modprobe -r r8619"
