# etc/udev/rules.d/redGearCosmos71_autoeffects/87-disable-requiredb.rules
#SUBSYSTEM=="net", SUBSYSTEMS=="pci", ATTRS{idVendor}=="14e4", RUN+="/usr/bin/modprobe -r wl0"
