# etc/profile.d/export.sh

export PATH="~/.local/bin:$PATH" 
export BAT="BAT1"
export AC="ACAD"
