# etc/udev/rules.d/99-t2u-led-off.rules
ACTION=="add", SUBSYSTEM=="leds", \
  ATTRS{idVendor}=="2357", ATTRS{idProduct}=="011e", \
  RUN+="/bin/sh -c 'echo none > /sys/class/leds/%k/trigger'"



